# .docker/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:php-fpm]
command=/usr/local/sbin/php-fpm
stdout_logfile=/dev/stdout
